# deploy/systemd/airbridge.service
[Unit]
Description=AirBridge Core Service (AirPlay -> Alexa)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=airbridge
Group=airbridge
WorkingDirectory=/opt/airbridge
EnvironmentFile=/etc/airbridge/airbridge.env
ExecStart=/usr/bin/node /opt/airbridge/dist/src/index.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/airbridge /run/airbridge /etc/airbridge /etc/credstore.encrypted
RuntimeDirectory=airbridge
RuntimeDirectoryMode=0750
UMask=0027

[Install]
WantedBy=multi-user.target
