# deploy/systemd/cloudflared-airbridge.service
[Unit]
Description=Cloudflared tunnel for AirBridge stream and skill endpoints
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nobody
Group=nogroup
ExecStart=/usr/bin/cloudflared tunnel --config /etc/airbridge/cloudflared.yml run
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/airbridge /etc/cloudflared

[Install]
WantedBy=multi-user.target
